# jitsi/etc/systemd/system/certbot.service.d/override.conf
[Service]
ExecStartPost=/usr/local/bin/certbot-post.sh
